# falcon-swcollector.conf
[program:falcon-swcollec]
command=/usr/local/open-falcon/falcon-swcollector/control run
autostart = true
startsecs = 5
autorestart=true
user=root
