# octopux-swcollector.conf
[program:swcollector]
command=/usr/local/octopux/octopux-swcollector/control run
autostart = true
startsecs = 5
autorestart=true
user=root
